# infra/systemd/caddy-self-heal.service
[Unit]
Description=uTrade Caddy self-heal check
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/caddy-self-heal.sh
